# packaging/cm.service
[Unit]
Description=Config Manager — headless node management
Documentation=https://github.com/msutara/config-manager-core
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/cm --headless
WorkingDirectory=/var/lib/cm
Restart=on-failure
RestartSec=5

# Environment defaults for systemd.
# Remove or comment out individual lines to let /etc/cm/config.yaml take effect.
# Environment=CM_LOG_LEVEL=info

# Logging — use systemd journal; cm also writes to /var/log/cm/cm.log
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cm

# Hardening — relaxed because CM manages system packages and network config.
ProtectHome=yes
ProtectControlGroups=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
